# deploy/camjam/systemd/camjam-control.service
[Unit]
Description=CamJam Robot Control Service
After=network-online.target pigpiod.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/camjam
ExecStart=/opt/camjam/bin/control --config /etc/camjam/pins.yaml --servo-calibration /var/lib/camjam/servos.yaml
Environment=CAMJAM_PIN_CONFIG=/etc/camjam/pins.yaml
Environment=CAMJAM_SERVO_CALIBRATION=/var/lib/camjam/servos.yaml
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
