# deploy/camjam/systemd/camjam-tls-renew.service
[Unit]
Description=Renew CamJam TLS Certificates

[Service]
Type=oneshot
User=root
ExecStart=/opt/camjam/bin/tls_renew.sh
